# config/supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
logfile_maxbytes=10MB
pidfile=/tmp/supervisord.pid
user=root

[program:xvfb]
command=/usr/bin/Xvfb :1 -screen 0 1280x800x24 -ac +extension GLX +render -noreset
autorestart=true
priority=10
stdout_logfile=/tmp/xvfb.log
stderr_logfile=/tmp/xvfb.log

[program:fluxbox]
command=/usr/bin/fluxbox
environment=DISPLAY=":1",HOME="/root"
autorestart=true
priority=20
startsecs=2
stdout_logfile=/tmp/fluxbox.log
stderr_logfile=/tmp/fluxbox.log

[program:x11vnc]
command=/usr/bin/x11vnc -display :1 -nopw -listen localhost -xkb -forever -shared -ncache 10 -quiet
autorestart=true
priority=30
startsecs=3
stdout_logfile=/tmp/x11vnc.log
stderr_logfile=/tmp/x11vnc.log

[program:novnc]
command=/usr/bin/websockify --web=/usr/share/novnc 8090 localhost:5900
autorestart=true
priority=40
startsecs=4
stdout_logfile=/tmp/novnc.log
stderr_logfile=/tmp/novnc.log

[program:rcommander]
command=/usr/bin/Rscript -e "library(Rcmdr)"
environment=DISPLAY=":1",HOME="/root"
directory=/data
autorestart=false
priority=50
startsecs=5
stdout_logfile=/tmp/rcommander.log
stderr_logfile=/tmp/rcommander.log
